# reflect30.service
[Unit]
Description=Reflect-30 Telegram Bot
After=network.target

[Service]
# Modify these to match your user and installation path
User=ubuntu
WorkingDirectory=/home/ubuntu/Reflectbot
ExecStart=/usr/bin/python3 /home/ubuntu/Reflectbot/bot.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/Reflectbot/.env

[Install]
WantedBy=multi-user.target
